# systemd/scale-volume@.service
[Unit]
Description=Scale boot volume

[Service]
Type=oneshot
EnvironmentFile=/etc/scale-volumes/%i.conf
ExecStart=/usr/lib/scale-volume
